# deployment/oaipmh.conf
[program:oaipmh]
command=/home/mark/jper-oaipmh/bin/gunicorn -c /home/mark/jper-oaipmh/src/jper-oaipmh/deployment/gconf.py service.web:app
user=mark
directory=/home/mark/jper-oaipmh/src/jper-oaipmh
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
